# infrastructure/docker/migration-day4/scripts/day19_read_switch_hourly.cron.example
# Run Day 19 read-switch snapshot every hour on minute 40 (UTC)
# Install: crontab infrastructure/docker/migration-day4/scripts/day19_read_switch_hourly.cron.example

SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

40 * * * * cd /workspaces/Data-Platform && bash infrastructure/docker/migration-day4/scripts/run_day19_read_switch_hourly.sh >> documentation/migration/reports/DAY19_READ_SWITCH_CRON.log 2>&1
