# infrastructure/docker/migration-day4/scripts/day20_traffic_hourly.cron.example
# Run Day 20 traffic snapshot every hour on minute 45 (UTC)
# Install: crontab infrastructure/docker/migration-day4/scripts/day20_traffic_hourly.cron.example

SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

45 * * * * cd /workspaces/Data-Platform && bash infrastructure/docker/migration-day4/scripts/run_day20_traffic_hourly.sh >> documentation/migration/reports/DAY20_TRAFFIC_CRON.log 2>&1
